# system/files/rc.local
#!/bin/sh -e

##################################
# **** DO NOT EDIT THIS FILE ****
# This file is managed by Salt
# Any changes will be overwritten
##################################

echo never > /sys/kernel/mm/transparent_hugepage/enabled
exit 0
